# graphite/service/graphite-web/run
#!/bin/bash
exec 2>&1
PYTHONPATH=/opt/graphite/webapp gunicorn -b 0.0.0.0:8000 -w 2 --log-level debug graphite.wsgi
