# graphite/service/carbon/run
#!/bin/bash
exec 2>&1
/opt/graphite/bin/carbon-cache.py --config /opt/graphite/conf/carbon.conf --debug start
